# ansible-security/host_final/spark_log_ingest/log_ingest.service
[Unit]
Description=Processar logs Suricata com PySpark

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/scripts/process_suricata_logs.py
